# deploy/systemd/rpg-alert-emitter.timer
[Unit]
Description=Run RPG alert emitter every 60 seconds

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
Unit=rpg-alert-emitter.service
AccuracySec=1s

[Install]
WantedBy=timers.target
