# deploy/systemd/rpg-alert-emitter.service
[Unit]
Description=RPG observability alert emitter
After=network-online.target rpg-backend.service
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/rpg-backend
EnvironmentFile=/opt/rpg-backend/.env
ExecStart=/bin/sh -lc '/opt/rpg-backend/.venv/bin/python scripts/emit_runtime_alerts.py --window-seconds "${APP_OBS_ALERT_WINDOW_SECONDS:-300}"'
